# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:websocket-server]
command=/usr/local/bin/php /var/www/modules/websocket/server.php
autostart=true
autorestart=true
user=www-data
directory=/var/www
